# TP 3-LINUX/yt.service
[Unit]
Description=Youtube downloader

[Service]
ExecStart=sudo bash /home/zmehdi/srv/yt/yt-v2.sh

[Install]
WantedBy=multi-user.target